# etc/service/powerdns/run
#!/bin/bash -e

exec /usr/sbin/pdns_server 2>&1
